# backend/scripts/crontab
# Run email digest agent at 8:00 AM from Monday to Friday (PST)
0 8 * * * cd /app && /usr/local/bin/uv run main.py >> /var/log/emailagent/cron.log 2>&1
